# device/lights-collector.service
[Unit]
Description=Lights Data Collector Service
After=network-online.target
Wants=network-online.target

[Service]
# Change 'pi' to your username if different (e.g., 'dhooper')
User=dhooper
Group=dhooper

# Set the working directory so relative paths (logs, state file) work
WorkingDirectory=/home/dhooper/gemini/lights/device

# Use the Python executable from your virtual environment
ExecStart=/home/dhooper/gemini/lights/device/venv/bin/python collect.py

# Restart the service if it crashes
Restart=always
RestartSec=10

# Output logs to system journal (view with journalctl -u lights-collector)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
